# src/fty-sensor-gpio.service.in
[Unit]
Description=fty-sensor-gpio service
Requires=network.target malamute.service fty-asset.service
After=network.target malamute.service fty-asset.service
PartOf=bios.target

[Service]
Type=simple
User=bios
Restart=always
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
Environment='SYSTEMD_UNIT_FULLNAME=%n'
ExecStart=@prefix@/bin/fty-sensor-gpio -c @sysconfdir@/@PACKAGE@/fty-sensor-gpio.cfg
ExecStop=/bin/kill -9 $MAINPID
Restart=always

[Install]
WantedBy=bios.target
